# supervisord.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/usr/sbin/varnishd -F -a %(ENV_VARNISH_LISTEN_ADDRESS)s:%(ENV_VARNISH_LISTEN_PORT)s -f %(ENV_VARNISH_VCL_CONF)s -T %(ENV_VARNISH_ADMIN_LISTEN_ADDRESS)s:%(ENV_VARNISH_ADMIN_LISTEN_PORT)s -t %(ENV_VARNISH_TTL)s -p thread_pool_min=%(ENV_VARNISH_MIN_THREADS)s -p thread_pool_max=%(ENV_VARNISH_MAX_THREADS)s -p thread_pool_timeout=%(ENV_VARNISH_THREAD_TIMEOUT)s -S %(ENV_VARNISH_SECRET_FILE)s -s %(ENV_VARNISH_STORAGE)s
autorestart=true

[program:varnishncsa]
command=/usr/bin/varnishncsa -F "%(ENV_VARNISH_NCSA_LOGFORMAT)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=10
